# pixelfed-setup.service
[Unit]
Description=Setup Pixelfed

[Service]
WorkingDirectory=/home/pixelfed/pixelfed
Type=oneshot
ExecStart=/home/pixelfed/setup.sh

[Install]
WantedBy=default.target